# archive/com.bartturczynski.ccfix-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bartturczynski.ccfix-watch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/bartturczynski/Projects/claude-copy-fix/ccfix.py</string>
        <string>--watch</string>
    </array>

    <!-- Start at login and restart if it ever exits. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Logs, for if you ever want to see what it did / debug. -->
    <key>StandardOutPath</key>
    <string>/tmp/ccfix-watch.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ccfix-watch.log</string>
</dict>
</plist>
